# run.ini
[program:benchmark]
directory=/home/isu-user/isucon/bench
command=/home/isu-user/isucon/env.sh carton exec -- ./server.pl
user=isu-user
stdout_logfile=/tmp/isucon.bench.log
stderr_logfile=/tmp/isucon.bench.log
autostart=true

[program:isucon_perl]
directory=/home/isu-user/isucon/webapp/perl
command=/home/isu-user/isucon/env.sh carton exec -- plackup -s Starman --workers 10 --port 5000 -E production --disable-keepalive app.psgi
user=isu-user
stdout_logfile=/tmp/isucon.perl.log
stderr_logfile=/tmp/isucon.perl.log
autostart=false

[program:isucon_node]
directory=/home/isu-user/isucon/webapp/nodejs
command=/home/isu-user/isucon/env.sh node server.js
user=isu-user
stdout_logfile=/tmp/isucon.node.log
stderr_logfile=/tmp/isucon.node.log
stopsignal=INT
autostart=false

[program:isucon_ruby]
directory=/home/isu-user/isucon/webapp/ruby
command=/home/isu-user/isucon/env.sh bundle exec foreman start
user=isu-user
stdout_logfile=/tmp/isucon.ruby.log
stderr_logfile=/tmp/isucon.ruby.log
autostart=true

[program:isucon_python]
directory=/home/isu-user/isucon/webapp/python
command=/home/isu-user/isucon/env.sh gunicorn -c gunicorn_config.py -w 10 app:app
user=isu-user
stdout_logfile=/tmp/isucon.python.log
stderr_logfile=/tmp/isucon.python.log
autostart=false

[program:isucon_go]
directory=/home/isu-user/isucon/webapp/go
command=/home/isu-user/isucon/env.sh ./app
user=isu-user
stdout_logfile=/tmp/isucon.go.log
stderr_logfile=/tmp/isucon.go.log
autostart=false
